# Development/Packer/site-cookbooks/environment/templates/default/mailcatcher.upstart.conf.erb
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN

description		"Mailcatcher SMTP dev tool"
author			"Josef Florian Glatz <jousch@gmail.com>"

start on startup
stop on shutdown

exec /usr/local/bin/mailcatcher --http-ip <%= node['mailcatcher']['http-ip'] %> --http-port <%= node['mailcatcher']['http-port'] %> --smtp-port <%= node['mailcatcher']['smtp-port'] %>